# infra/systemd/lab-management-backend@.service
[Unit]
Description=Lab Management backend for %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i/Documents/HelloWorld/backend
EnvironmentFile=/home/%i/Documents/HelloWorld/.env.rpi
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/%i/Documents/HelloWorld/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
